# resize_rootfs.service
[Unit]
Description=resize_rootfs
 
[Service]
Type=oneshot
ExecStart=/bin/bash /bin/resize_rootfs.sh

[Install]
WantedBy=multi-user.target
